# files/etc/hotplug.d/net/55-rps-rx-tuning
#!/bin/sh

# Spread RX softirq across CPUs 1-3 for LAN/WAN interfaces.
# airoha_eth driver has no RSS, so without this all multi-flow RX
# pins to CPU0 (the hard IRQ CPU). Mask 'e' = 1110 = CPUs 1,2,3,
# leaving CPU0 free for hard IRQ handling.

[ "$ACTION" = "add" ] || exit 0

DEV="${DEVICENAME:-$INTERFACE}"
[ -n "$DEV" ] || exit 0

case "$DEV" in
	lan*|wan*) ;;
	*) exit 0 ;;
esac

for q in /sys/class/net/"$DEV"/queues/rx-*/rps_cpus; do
	[ -w "$q" ] && echo e > "$q"
done

for q in /sys/class/net/"$DEV"/queues/rx-*/rps_flow_cnt; do
	[ -w "$q" ] && echo 4096 > "$q"
done

exit 0
